# cron.jobs/cron.daily/daily_jobs
#!/bin/bash

# Updates 'locate' database
updatedb

# Delete conflicted copies
find /home/bryan/Dropbox -name "*aphrodite's conflicted copy*" -delete
