# cron.jobs/cron.daily/daily_jobs
#!/bin/bash

# Updates 'locate' database
updatedb

# Delete conflicted copies
find /home/bryan/Dropbox -name "*'s conflicted copy*" -delete

# Clear /home/bryan/Dropbox/tmp
rm -rf /home/bryan/Dropbox/tmp
mkdir /home/bryan/Dropbox/tmp
