# cron.jobs/cron.daily/daily_jobs
#!/bin/bash

/usr/bin/updatedb
/usr/local/bin/etc-generator
/usr/local/bin/clean_filesystem
/usr/local/bin/zopen -C -q
/usr/local/bin/merge_rocketbook_notes
