# cron.jobs/cron.daily/daily_jobs
#!/bin/bash

set -e

/usr/bin/updatedb
clean_filesystem
zopen -C -q
mopen -C
check_shv_logs

# zopen and mopen create log files owned by root.
chown bryan:bryan /var/tmp/*.log
