# cron.jobs/cron.daily/daily_jobs
#!/bin/bash

set -e

/usr/bin/updatedb
clean_filesystem
zopen -C -q
mopen -C
check_shv_logs

# Run tests for my local scripts.
run_all_script_tests

# Generate daily summary of process accounting (needed for `sar`).
/usr/lib64/sa/sa2 -A &

# zopen and mopen create log files owned by root.
chown bryan:bryan /var/tmp/*.log
