# cron.jobs/cron.weekly/weekly_jobs
#!/bin/bash

set -e

sort_qb_quickmarks
clean_recent_file_log
clean_mpv_history
cache_cheaders
