# cron.jobs/cron.weekly/weekly_backup
#!/bin/bash

source /home/bryan/Sync/bin/cron.jobs/backup.sh

function main() {
    ETBB=$((60 * 60 * 24 * 7)) # seconds in a week

    backup /etc etc/weekly "${DEFAULT_R}"
    backup_home weekly "${DEFAULT_R}"

    post_backup_hook
}

main "$@"
