# cron.jobs/cron.weekly/weekly_jobs
#!/bin/bash
/home/bryan/Dropbox/scripts/main/textme "cron.weekly/weekly_jobs: This weekly test message is necessary in order to ensure this phone number stays active."

rm -f /home/bryan/Dropbox/logs/aphrodite-motion/*
