# bin/cron.jobs/cron.daily/daily_backup
#!/bin/bash

# dropbox backup
DROPBOX='/media/bryan/hercules/backup/dropbox/daily'
rsync -av --delete /home/bryan/Dropbox/ $DROPBOX
date > $DROPBOX/backup.txt

# etc backup
ETC='/media/bryan/hercules/backup/etc/daily'
rsync -av --delete /etc $ETC
date > $ETC/backup.txt

# world file backup
cp /var/lib/portage/world /media/bryan/hercules/backup/portage/daily.txt
