# bin/cron.jobs/cron.daily/daily_jobs
#!/bin/bash

# Checks if Maintenance id Due
epmaint --check

# Updates 'locate' database
updatedb

# Delete conflicted copies
find /home/bryan/Dropbox -name "*'s conflicted copy*" -delete
