# bin/cron.jobs/cron.daily/daily_backup
#!/bin/bash

# dropbox backup
DROPBOX='/media/bryan/hercules/backup/dropbox/daily'
rsync -av --delete /home/bryan/Dropbox/ $DROPBOX
date > $DROPBOX/backup.txt

# package list backup
/usr/local/bin/eplist -o > /media/bryan/hercules/backup/portage/daily.txt
