# cron.jobs/cron.monthly/monthly_backup
#!/bin/bash

source /home/bryan/Sync/bin/cron.jobs/backup.sh

function main() {
    ETBB=$((60 * 60 * 24 * 30)) # seconds in a month

    backup /etc etc/monthly "${DEFAULT_R}"
    backup_home monthly "${DEFAULT_R}"

    # Yearly Backup
    if [[ $(date +%m) -eq 1 ]]; then
        ETBB=$((60 * 60 * 24 * 365)) # seconds in a year

        backup /etc etc/yearly "${DEFAULT_R}"
        backup_home yearly "${DEFAULT_R}"
    fi

    post_backup_hook
}

main "$@"
